# scripts/systemd/restream-nodecg@.service
[Unit]
Description=Restreaming Automation NodeCG Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/opt/restreaming_automation/nodecg
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/restreaming_automation
PrivateTmp=true

[Install]
WantedBy=multi-user.target
